# scripts/crontab.cron
0 12 * * * ~josh/proj/nyx/scripts/packages/refresh-packages-cron.sh
