# scripts/refresh-packages.cron
0 12 * * * ~josh/proj/nyx/scripts/refresh-packages-cron.sh
